# conf.d/x11vnc.conf
[program:x11vnc]
command=x11vnc -forever -shared -auth guess -frostherovnc
autorestart=true
